# crontab
5 * * * * root cd miner && ./cpu.sh

@reboot root cd miner && ./cpu.sh
